# rootfs/etc/init.d/initpywrk
#!/sbin/openrc-run

command=/usr/bin/init_worker
